# rpi/etc/rc.local
#!/bin/bash

set -e

for script in /etc/rc.local.d/*; do
  $script
done

exit 0
